# crontab
# Run crossword download+upload every day at 01:00
0 1 * * * /app/cron-entry.sh >> /var/log/rm-xword.log 2>&1
# (blank line required by cron)
